# files/system/etc/systemd/system/nix-install.service
[Unit]
Description=Install Nix Package Manager (first boot)
ConditionPathExists=!/var/nix/.installed
After=network-online.target nix.mount
Wants=network-online.target
Requires=nix.mount

[Service]
Type=oneshot
ExecStart=/usr/bin/bash -c 'until curl -fsSL https://install.determinate.systems/nix -o /tmp/nix-installer.sh; do echo "Waiting for network..."; sleep 5; done && sh /tmp/nix-installer.sh install linux --init none --no-confirm --extra-conf "sandbox = false" --extra-conf "filter-syscalls = false" && touch /var/nix/.installed'
RemainAfterExit=yes
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
